# roles/proxy/templates/nginx-health-check.service.j2
# {{ ansible_managed }}
[Unit]
Description=Nginx Upstream Health Check
After=nginx.service

[Service]
Type=oneshot
ExecStart=/usr/local/bin/nginx-health-check.sh
